# bradbot.service
[Unit]
Description=BradBot Discord Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bradbot
Environment=PATH=/home/ubuntu/bradbot/venv/bin
Environment=DB_HOST=ontgcoj4pznoqbq4aznk64mehy.dsql.us-east-1.on.aws
Environment=DB_PORT=5432
Environment=DB_NAME=postgres
Environment=DB_USER=BradBotRole
Environment=USE_IAM_AUTH=true
Environment=AWS_REGION=us-east-1
ExecStart=/home/ubuntu/bradbot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
