# bradbot.service
[Unit]
Description=BradBot Discord Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/bradbot
Environment=PATH=/home/ubuntu/bradbot/venv/bin
ExecStart=/home/ubuntu/bradbot/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
